# SOURCE_CODE/deploy/systemd/nossh-server.service
[Unit]
Description=nossh bridge server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nossh-server -config /etc/nossh/server.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
